# root/etc/services.d/web/run
#!/usr/bin/with-contenv bash

cd /gclient || exit
HOME="/config" exec node app.js